# debian/postrm
#!/bin/sh
# Removed set -e to prevent script from failing when commands fail

cleanup_files() {
    echo "Cleaning up Self service files..."
    
    rm -f /var/run/self.pid
    
    rm -rf /usr/lib/self
    
    rm -f /usr/sbin/self-start
    rm -f /usr/sbin/self-stop
    rm -f /usr/sbin/self-tool
    
    if [ -d "/sys/fs/bpf/self" ]; then
        rm -rf /sys/fs/bpf/self
    fi
    
    echo "File cleanup completed."
}

case "$1" in
    remove)
        if [ -x "$(command -v systemctl)" ]; then
            echo "Disabling Self service..."
            systemctl disable self.service || echo "Warning: Could not disable service"
            rm -f /etc/systemd/system/self.service
            systemctl daemon-reload
        fi
        
        cleanup_files
    ;;
    
    purge)
        if [ -f "/etc/systemd/system/self.service" ]; then
            rm -f /etc/systemd/system/self.service
            if [ -x "$(command -v systemctl)" ]; then
                systemctl daemon-reload
            fi
        fi
        
        cleanup_files
        
        rm -rf /etc/self
    ;;
    
    upgrade|failed-upgrade|abort-install|abort-upgrade|disappear)
        echo "No cleanup needed for action: $1"
    ;;

    *)
        echo "postrm does not recognize argument '$1'" >&2
        exit 1
    ;;
esac

exit 0 