# debian/self.service
[Unit]
Description=Self-learning firewall based on eBPF
After=network.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/sbin/self-start
ExecStop=/usr/sbin/self-stop
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target